# scripts/invosync-backup-task.xml
<?xml version="1.0" encoding="UTF-16"?>
<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Description>InvoSync hourly MongoDB backup</Description>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <Repetition>
        <Interval>PT1H</Interval>
        <Duration>P1D</Duration>
        <StopAtDurationEnd>false</StopAtDurationEnd>
      </Repetition>
      <StartBoundary>2026-01-01T00:00:00</StartBoundary>
      <Enabled>true</Enabled>
    </CalendarTrigger>
  </Triggers>
  <Actions Context="Author">
    <Exec>
      <Command>python</Command>
      <Arguments>C:\Users\Admin\Desktop\Project-Pauldirac\scripts\backup_schedule.py --quiet</Arguments>
      <WorkingDirectory>C:\Users\Admin\Desktop\Project-Pauldirac</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
